# scheduledJobs/mysqlbackup.service
[Unit]
Description=MySQL Backup Job

[Service]
Type=simple
Restart=no
ExecStart=/usr/local/bin/mysqlbackup.sh